# backup/.crontab
0 0 * * * python3 /home/[username]/backup_db.py
